# .config/launchd/com.4jp.pipeline.daily-deferred.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.4jp.pipeline.daily-deferred</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/4jp/Workspace/4444J99/application-pipeline/.venv/bin/python</string>
        <string>/Users/4jp/Workspace/4444J99/application-pipeline/scripts/check_deferred.py</string>
        <string>--alert</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/4jp/Workspace/4444J99/application-pipeline</string>
    <key>StandardOutPath</key>
    <string>/Users/4jp/System/Logs/pipeline-daily-deferred.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/4jp/System/Logs/pipeline-daily-deferred.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
